# x/udev_rules/90-usb-device-plugin.rules
TTRS{idVendor}=="152d", ATTRS{idProduct}=="2329", RUN+="xset r rate 300 80"
