# de/udev-rules/90-usb-device-plugin.rules
ACTION=="add", SUBSYSTEM=="usb_device", RUN+="xset r rate 300 80"
